# rqscheduler.conf
[program:rqscheduler]
command=python manage.py rqscheduler
directory=/vagrant/progr/py/rq/django-test-rq
environment=PATH="/home/vagrant/.virtualenvs/rq/bin"
user=vagrant
redirect_stderr=true
